# src/icsync/templates/com.icsync.import.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{label}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{python}</string>
        <string>-m</string>
        <string>icsync</string>
        <string>run</string>
    </array>

    <key>StartInterval</key>
    <integer>{interval}</integer>
    <key>WatchPaths</key>
    <array>
        <string>{inbox}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{stdout}</string>
    <key>StandardErrorPath</key>
    <string>{stderr}</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
